# deploy/systemd/faceseek-frontend.service
[Unit]
Description=FaceSeek Frontend (Next.js)
After=network.target faceseek-backend.service

[Service]
Type=simple
WorkingDirectory=/opt/faceseek/repo/frontend
EnvironmentFile=/etc/faceseek/frontend.env
User=faceseek
Group=faceseek
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/env npm run start -- -p 3000
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
